# deploy/systemd/chess-engine-dev.service
[Unit]
Description=Chess Engine backend (development, hot reload)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as root so it has permission to read ~/srv/
User=root
Group=root

# Absolute path to your dev backend codebase
WorkingDirectory=/root/srv/chess-engine/dev/backend

# Point systemd directly to your dev .env file
EnvironmentFile=/root/srv/chess-engine/dev/backend/.env

# Absolute path to Bun. We add --hot to automatically reload on changes.
ExecStart=/root/.bun/bin/bun --hot src/server.js

Restart=on-failure
RestartSec=2
TimeoutStopSec=30
KillSignal=SIGTERM
SyslogIdentifier=chess-engine-dev

[Install]
WantedBy=multi-user.target