# deploy/systemd/chess-engine-prod.service
[Unit]
Description=Chess Engine backend (production)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as root so it has permission to read ~/srv/
User=root
Group=root

# Absolute path to your backend codebase
WorkingDirectory=/root/srv/chess-engine/prod/backend

# Point systemd directly to your .env file
EnvironmentFile=/root/srv/chess-engine/prod/backend/.env

# Absolute path to Bun (no symlinks needed)
ExecStart=/root/.bun/bin/bun run src/server.js

# Graceful shutdown and restart logic
Restart=on-failure
RestartSec=2
TimeoutStopSec=30
KillSignal=SIGTERM
SyslogIdentifier=chess-engine

[Install]
WantedBy=multi-user.target